# configs/tide_update.timer
[Unit]
Description=Get today's tide data

[Timer]
OnCalendar=*-*-* 01:30:00
Persistent=True
Unit=tide_update.service

[Install]
WantedBy=timers.target
